# systemd/proxmox-review-sync.service
[Unit]
Description=Publish sanitized Proxmox review data
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/flock -n /run/proxmox-review-sync.lock /usr/local/bin/proxmox-review-sync.sh --env @ENV_PATH@
